# pi/systemd/sr33-archiver.service
[Unit]
Description=SR33 full-resolution local archive (Signal K -> onboard SQLite)
After=network-online.target signalk.service
Wants=network-online.target

[Service]
Type=simple
# Drop real values in /etc/sr33/archiver.env on the Pi (BOAT_ID, ARCHIVE_DB, SIGNALK_WS, ...)
EnvironmentFile=/etc/sr33/archiver.env
ExecStart=/usr/bin/python3 /opt/sr33/pi/archiver/archiver.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
