# pi/systemd/sr33-uplink.service
[Unit]
Description=SR33 telemetry uplink (Signal K -> VPS ingestion)
After=network-online.target signalk.service
Wants=network-online.target

[Service]
Type=simple
# Drop real values in /etc/sr33/uplink.env on the Pi (CAN_IFACE=can0, VPS_URL, INGEST_TOKEN, ...)
EnvironmentFile=/etc/sr33/uplink.env
ExecStart=/usr/bin/python3 /opt/sr33/pi/uplink/uplink.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
